# roles/clients/templates/restic-backup.service.j2
[Unit]
Description=Restic Backup Service
Documentation=man:restic(1)
Documentation=https://restic.readthedocs.io/en/stable/
Requires=local-fs.target
Requires=network.target

[Service]
Type=oneshot
ExecStart={{ restic_backup_client_dir }}/make_backup.sh
WorkingDirectory={{ restic_backup_client_dir }}

## Security hardening: see systemd.exec(7) and systemd.resource-control(5)
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths={{ restic_backup_client_cache_dir }}
ReadWritePaths=/var/spool/postfix/maildrop

UMask=077

ProtectHostname=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true

ProtectProc=invisible
ProcSubset=pid

PrivateDevices=true
DevicePolicy=closed
MemoryDenyWriteExecute=true

NoNewPrivileges=true
SystemCallArchitectures=native
LockPersonality=true

RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

CapabilityBoundingSet=~CAP_SYS_PACCT
CapabilityBoundingSet=~CAP_KILL
#CapabilityBoundingSet=~CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER CAP_IPC_OWNER # needed
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
CapabilityBoundingSet=~CAP_IPC_LOCK
CapabilityBoundingSet=~CAP_SYS_TTY_CONFIG
CapabilityBoundingSet=~CAP_SYS_BOOT
CapabilityBoundingSet=~CAP_SYS_CHROOT
CapabilityBoundingSet=~CAP_BLOCK_SUSPEND
CapabilityBoundingSet=~CAP_LEASE
CapabilityBoundingSet=~CAP_CHOWN CAP_FSETID CAP_SETFCAP
CapabilityBoundingSet=~CAP_SETUID CAP_SETGID CAP_SETPCAP
CapabilityBoundingSet=~CAP_MAC_ADMIN CAP_MAC_OVERRIDE
CapabilityBoundingSet=~CAP_SYS_PTRACE
CapabilityBoundingSet=~CAP_SYS_NICE CAP_SYS_RESOURCE
CapabilityBoundingSet=~CAP_NET_ADMIN
CapabilityBoundingSet=~CAP_NET_BIND_SERVICE CAP_NET_BROADCAST CAP_NET_RAW
CapabilityBoundingSet=~CAP_AUDIT_CONTROL CAP_AUDIT_READ CAP_AUDIT_WRITE

SystemCallFilter=~@clock
SystemCallFilter=~@cpu-emulation
SystemCallFilter=~@debug
SystemCallFilter=~@module
SystemCallFilter=~@mount
SystemCallFilter=~@obsolete
#SystemCallFilter=~@privileged # needed
SystemCallFilter=~@raw-io
SystemCallFilter=~@reboot
#SystemCallFilter=~@resources # needed
SystemCallFilter=~@swap
